# deploy/longtail-ssh-firewall.service
[Unit]
Description=Restrict SSH access to WireGuard network
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/longtail
Environment=NET_TOOLS_WG_IFACE=wg0
Environment=NET_TOOLS_WG_CIDR=10.66.0.0/24
Environment=NET_TOOLS_SSH_PORT=22
EnvironmentFile=-/etc/longtail/longtail.env
ExecStart=/bin/bash /opt/longtail/deploy/apply-ssh-access.sh
RemainAfterExit=yes
User=root
Group=root

[Install]
WantedBy=multi-user.target
